# updateAngle.service
[Unit]
Description=Autoupdate angle

[Service]
ExecStart=python updateAngle.py
WorkingDirectory=/root

[Install]
WantedBy=multi-user.target
